# bowei_ai_dashboard/scripts/bowei.service
[Unit]
Description=Bowei AI Dashboard backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/bowei_ai_dashboard
Environment=APP_ENV=production
EnvironmentFile=-/etc/bowei-ai-dashboard.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/bowei_ai_dashboard/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
